# crontab
# At every 10th minute.
* * * * * /usr/bin/python3 /root/main.py > /proc/1/fd/1 2>/proc/1/fd/2
